# maestro@.service
[Unit]
Description=Maestro agent orchestrator (%i)
After=network.target

[Service]
Type=simple
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin:%h/.bun/bin:%h/.npm-global/bin:%h/go/bin
ExecStart=/usr/local/bin/maestro run --config %h/.maestro/maestro-%i.yaml
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
